# ro-grf-analysis-tools.timer
[Unit]
Description=Timer for GRF Analysis tools

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=false

[Install]
WantedBy=timers.target
